# crontab
*/1 * * * * /bin/sh /usr/local/bin/s3-sync.sh
